# etc/systemd/phoned-ringerd.service.in
[Unit]
Description=phoned ringer service
After=phoned-modemd.service

[Service]
Type=simple
Restart=on-failure
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/phoned-ringerd

[Install]
WantedBy=multi-user.target