# etc/systemd/phoned-audiod.service.in
[Unit]
Description=phoned audio service
After=phoned-modemd.service phoned-handsetd.service

[Service]
Type=simple
Restart=on-failure
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/phoned-audiod

[Install]
WantedBy=multi-user.target