# etc/systemd/phoned-diald.service.in
[Unit]
Description=phoned dial service

[Service]
Type=simple
Restart=on-failure
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/phoned-diald

[Install]
WantedBy=multi-user.target